# scripts/com.igagent.v25.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.igagent.v25.watchdog</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>__IG_AGENT_ROOT__/scripts/watchdog.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>__IG_AGENT_ROOT__</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>__IG_AGENT_ROOT__/src/data/logs/watchdog_launchd.log</string>
	<key>StandardErrorPath</key>
	<string>__IG_AGENT_ROOT__/src/data/logs/watchdog_launchd.log</string>
</dict>
</plist>
