# scripts/com.igagent.v25.profitability.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.igagent.v25.profitability</string>
	<key>ProgramArguments</key>
	<array>
		<string>__PYTHON_BIN__</string>
		<string>__IG_AGENT_ROOT__/scripts/profitability_report.py</string>
		<string>--reconcile</string>
		<string>--days</string>
		<string>14</string>
	</array>
	<key>WorkingDirectory</key>
	<string>__IG_AGENT_ROOT__</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>__IG_AGENT_ROOT__/src</string>
	</dict>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>0</integer>
		<key>Hour</key>
		<integer>8</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>__IG_AGENT_ROOT__/src/data/logs/profitability_report.log</string>
	<key>StandardErrorPath</key>
	<string>__IG_AGENT_ROOT__/src/data/logs/profitability_report.log</string>
</dict>
</plist>
