# scripts/com.igagent.v29roadmap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.igagent.v29roadmap</string>
	<key>ProgramArguments</key>
	<array>
		<string>__PYTHON_BIN__</string>
		<string>__IG_AGENT_ROOT__/scripts/roadmap_morning_report.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>__IG_AGENT_ROOT__</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>__IG_AGENT_ROOT__/src:__IG_AGENT_ROOT__/v26</string>
	</dict>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>5</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>__IG_AGENT_ROOT__/src/data/logs/v29_roadmap_morning.log</string>
	<key>StandardErrorPath</key>
	<string>__IG_AGENT_ROOT__/src/data/logs/v29_roadmap_morning.log</string>
</dict>
</plist>
